# deploy/ansible/databases/roles/podman/files/oracle.service.j2
[Unit]
Description=Oracle Database Free 26ai (podman)
After=network-online.target
Wants=network-online.target

[Service]
Restart=always
TimeoutStartSec=900
ExecStartPre=-/usr/bin/podman rm -f oracle
ExecStart=/usr/bin/podman run --name oracle \
  -p 1521:1521 \
  -e ORACLE_PWD={{ oracle_db_password }} \
  -v /data/oracle:/opt/oracle/oradata:Z \
  container-registry.oracle.com/database/free:latest-26ai
ExecStop=/usr/bin/podman stop -t 90 oracle

[Install]
WantedBy=multi-user.target
